# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vnc]
command=/usr/bin/vncserver -localhost no :1 -geometry 1920x1080 -depth 24
autostart=true
autorestart=true
stderr_logfile=/var/log/vnc.err.log
stdout_logfile=/var/log/vnc.out.log

[program:novnc]
command=/usr/bin/websockify -D --web=/usr/share/novnc 8080 localhost:5901
autostart=true
autorestart=true
stderr_logfile=/var/log/novnc.err.log
stdout_logfile=/var/log/novnc.out.log

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
stderr_logfile=/var/log/xrdp.err.log
stdout_logfile=/var/log/xrdp.out.log

[program:file-server]
command=/usr/bin/python3 /root/setup.py
autostart=true
autorestart=true
stderr_logfile=/var/log/fileserver.err.log
stdout_logfile=/var/log/fileserver.out.log